# contrib/systemd/nabsleep.timer
[Unit]
Description=Good Night

[Timer]
OnCalendar=Mon..Fri 20:00
OnCalendar=Sat,Sun 22:00
Persistent=true

[Install]
WantedBy=timers.target
