# contrib/systemd/nabsleep.service
[Unit]
Description=Sleep my rabbit

[Service]
Type=simple
ExecStart=/usr/local/bin/nabcli sleep
