# contrib/systemd/nabwakeup.service
[Unit]
Description=Wakeup my rabbit

[Service]
Type=simple
ExecStart=/usr/local/bin/nabcli wakeup
